# presearch.service
[Unit]
Description=Presearch Node

[Service]
ExecStart=/root/run

[Install]
WantedBy=multi-user.target
